# systemd/system_validator.service
# system_validator.service
# 豁｣蠑上ヱ繧ｹ・・root/System_Validator/APP_DIR/theaterverse_final/systemd/system_validator.service
# 逶ｮ逧・ｼ壹し繝ｼ繝薙せ襍ｷ蜍輔→繧ｻ繧ｭ繝･繝ｪ繝・ぅ繝上・繝峨ル繝ｳ繧ｰ・亥ｼｷ蛹問即・・
[Unit]
Description=System Validator (Theaterverse Final)
After=network-online.target
Wants=network-online.target
StartLimitBurst=3
StartLimitIntervalSec=60

[Service]
Type=simple
WorkingDirectory=/root/System_Validator/APP_DIR/theaterverse_final/current
ExecStart=/usr/bin/env python3 -m app.app_main
Restart=on-failure
RestartSec=3
EnvironmentFile=/root/System_Validator/APP_DIR/theaterverse_final/.env

# ---- Security Hardening ----
User=validator
Group=validator
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/root/System_Validator/APP_DIR/theaterverse_final/current/var
ReadOnlyPaths=/root/System_Validator/APP_DIR/theaterverse_final/current
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictSUIDSGID=true
RestrictRealtime=true
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallFilter=@system-service @basic-io @network-io ~@mount ~@privileged

# 繝昴・繝医ｒ菴ｿ縺・ｴ蜷医・萓具ｼ亥ｿ・ｦ√↓蠢懊§縺ｦ隱ｿ謨ｴ・・# ExecStartPre=/usr/bin/bash -lc 'test -S /run/system_validator.sock && rm -f /run/system_validator.sock || true'
# Socket=/run/system_validator.sock

[Install]
WantedBy=multi-user.target

# --- END OF STRUCTURE ---
# /root/System_Validator/APP_DIR/theaterverse_final/systemd/system_validator.service
# /root/System_Validator/APP_DIR/theaterverse_final/systemd/system_validator.service
# --- END OF STRUCTURE ---
